# templates/ubuntu/etc/init/margarine.conf.erb
description "Margarine Daemon: <%= @service %>"
author "Alex Brandt <alunduil@alunduil.com>"

setuid <%= node['margarine']['user'] %>
setgid <%= node['margarine']['group'] %>

console log

start on (
     started network-services
  or started networking
)

stop on (
     stopping network-services
  or stopping networking
)

exec /usr/local/bin/<%= @service %>
